# deploy/supervisor.conf
; /etc/supervisor/conf.d/fidelislogic.conf
;
; After copying:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl status fidelislogic-backend

[program:fidelislogic-backend]
directory=/app/backend
command=/app/backend/venv/bin/uvicorn server:app
    --host 127.0.0.1
    --port 8001
    --workers 1
    --proxy-headers
    --forwarded-allow-ips=127.0.0.1
    --log-level info
user=www-data
autostart=true
autorestart=true
startsecs=5
startretries=3
stopsignal=TERM
stopwaitsecs=15
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/fidelislogic-backend.out.log
stderr_logfile=/var/log/supervisor/fidelislogic-backend.err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
